# deploy/astrbot-launcher.service
[Unit]
Description=AstrBot Launcher Web Manager (Rust HTTP backend)
Documentation=https://github.com/AstrBotDevs/astrbot-launcher
After=network.target

[Service]
Type=simple
User=mulq
Group=mulq
WorkingDirectory=/opt/astrbot-launcher
ExecStart=/opt/astrbot-launcher/bin/astrbot_launcher
# 监听端口（代码默认 6190，可被 ASTRBOT_HTTP_PORT 覆盖）
Environment=ASTRBOT_HTTP_PORT=6190
# 数据目录：默认走 dirs crate（~/.local/share/astrbot-launcher）。
# 若需接管现有 /home/mulq/astrbot_launcher，待代码支持自定义 base-dir env 后再加 Environment。
# restart/stop_manager 走 systemd（见 commands.rs），Restart=on-failure 仅在进程异常退出时拉起。
Restart=on-failure
RestartSec=3
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
