# deploy/astrbot-launcher.service
[Unit]
Description=AstrBot Web Manager (Rust HTTP backend)
Documentation=https://github.com/ssc2991lyh/AstrBot_WebManager
After=network.target

[Service]
Type=simple
User=astrbot
Group=astrbot
# 系统用户无家目录，dirs crate 需要 HOME 才能解析数据目录，这里显式指定
WorkingDirectory=/var/lib/astrbot-launcher
Environment=ASTRBOT_HTTP_PORT=6190
Environment=HOME=/var/lib/astrbot-launcher
# restart/stop_manager 走 systemd（见 commands.rs），Restart=on-failure 仅在进程异常退出时拉起
ExecStart=/usr/local/bin/astrbot-launcher
Restart=on-failure
RestartSec=3
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
